# systemd/urouteros-backend.service
[Unit]
Description=URouterOS Backend Service
Documentation=https://github.com/dragonleehom/URouterPanel
After=network.target network-online.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/urouteros
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStart=/usr/bin/node /opt/urouteros/dist/index.js
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=urouteros-backend

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
